# luci-app-substore/root/etc/init.d/http-meta
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99
STOP=10

NAME=http-meta
PROG=/usr/bin/node
BUNDLE=/usr/libexec/substore/http-meta.bundle.js
TPL_SRC=/usr/libexec/substore/http-meta-tpl.yaml
# mihomo 内核由 +mihomo-meta 依赖保证安装时就存在，路径跟 Nikki、官方
# mihomo-meta/mihomo-alpha opkg 包完全一致，不再需要运行时探测/下载/提升。
SHARED_CORE=/usr/bin/mihomo

uci_get() {
	local ret=$(uci -q get substore.http_meta.$1)
	echo "${ret:-$2}"
}

start_service() {
	local enabled=$(uci_get enabled 0)
	[ "$enabled" = "1" ] || return 0

	if [ ! -f "$BUNDLE" ]; then
		echo "http-meta: bundle 文件缺失（$BUNDLE），请重新安装本包"
		return 0
	fi

	local meta_dir=$(uci_get meta_dir /etc/sub-store/http-meta)
	local external_core=$(uci_get external_core_path "")
	mkdir -p "$meta_dir/meta"

	# tpl.yaml 是构建期一起打包进包内的静态资源，首次启动时复制到运行目录，
	# 之后用户可自行修改，不会被后续启动覆盖。
	[ -f "$meta_dir/meta/tpl.yaml" ] || cp "$TPL_SRC" "$meta_dir/meta/tpl.yaml"

	# 内核：默认直接用 mihomo-meta 依赖装好的共享内核 /usr/bin/mihomo；
	# 如果用户在「外部内核路径」里手动指定了别的可执行文件，则优先使用它。
	if [ -n "$external_core" ] && [ -x "$external_core" ]; then
		ln -sf "$external_core" "$meta_dir/meta/http-meta"
	elif [ -x "$SHARED_CORE" ]; then
		ln -sf "$SHARED_CORE" "$meta_dir/meta/http-meta"
	else
		echo "http-meta: 未找到 mihomo 内核（$SHARED_CORE 不存在），mihomo-meta 依赖可能安装失败，请检查 opkg" >&2
		return 0
	fi

	local host=$(uci_get host 127.0.0.1)
	local port=$(uci_get port 9876)
	local authorization=$(uci_get authorization "")
	local max_port=$(uci_get max_available_port "")
	local min_port=$(uci_get min_available_port "")
	local disable_auto_clean=$(uci_get disable_auto_clean 0)
	local temp_folder=$(uci_get temp_folder "")
	local body_json_limit=$(uci_get body_json_limit 1mb)

	procd_open_instance
	procd_set_param command "$PROG" "$BUNDLE"
	procd_set_param respawn 3600 5 0
	procd_set_param stdout 1
	procd_set_param stderr 1

	procd_set_param env HOST="$host" PORT="$port" META_FOLDER="$meta_dir/meta" BODY_JSON_LIMIT="$body_json_limit"

	[ -n "$authorization" ] && procd_append_param env AUTHORIZATION="$authorization"
	[ -n "$max_port" ] && procd_append_param env META_MAX_AVAILABLE_PORT="$max_port"
	[ -n "$min_port" ] && procd_append_param env META_MIN_AVAILABLE_PORT="$min_port"
	[ "$disable_auto_clean" = "1" ] && procd_append_param env META_DISABLE_AUTO_CLEAN=true
	[ -n "$temp_folder" ] && procd_append_param env META_TEMP_FOLDER="$temp_folder"

	procd_close_instance
}

reload_service() {
	stop
	start
}

service_triggers() {
	procd_add_reload_trigger "substore"
}
